# com.tkr.mcp-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tkr.mcp-server</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/open</string>
    <string>-W</string>
    <string>-a</string>
    <string>/Users/plexariffic/@tkr-dev/tkr-swift-mcp/.build/release/TKRMCPServer.app</string>
    <string>--args</string>
    <string>--port</string>
    <string>42124</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/plexariffic/@tkr-dev/tkr-swift-mcp</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/plexariffic/Library/Logs/tkr/mcp-server.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/plexariffic/Library/Logs/tkr/mcp-server.stderr.log</string>

  <key>ThrottleInterval</key>
  <integer>5</integer>
</dict>
</plist>
